# templates/redmine.service
[Unit]
Description=Redmine
Documentation=http://redmine.org
Requires=postgresql.service
After=postgresql.service
Before=nginx.service

[Service]
Type=simple
WorkingDirectory=/opt/redmine/redmine-{{ redmine_version }}
User=redmine
Group=redmine
Environment=RAILS_ENV=production
ExecStart=/usr/bin/bundle exec puma -C /opt/redmine/redmine-{{ redmine_version }}/config/puma.rb

[Install]
WantedBy=multi-user.target

